# vagrant/provider-docker-varnish-monitoring-setup/ansible/roles/backend/templates/etc/systemd/system/demo_app.service.j2
[Unit]
Description=Demo App

[Service]
Type=simple
ExecStart=/usr/bin/java -jar -Dserver.port={{ demo_app_port }} /etc/demo_app/varnish-demo-app.jar 

SyslogIdentifier=demo_app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
